# .root-fs/etc/s6-overlay/s6-rc.d/prepare-filesystem/up
if { mkdir -p /var/lib/nginx/tmp }
if { mkdir -p /var/log/nginx }
chown -R www:www /app/storage
