# .root-fs/etc/s6-overlay/s6-rc.d/create-admin/up
with-contenv s6-setuidgid www php /app/artisan cattr:make:admin
