# .root-fs/etc/s6-overlay/s6-rc.d/prepare-cache/up
with-contenv s6-envuidgid www php82 /app/artisan config:cache
