# .root-fs/etc/s6-overlay/s6-rc.d/prepare-database/up
with-contenv s6-envuidgid www php /app/artisan migrate --force --seed --seeder=InitialSeeder
