# .root-fs/crontab
* * * * * php82 /app/artisan schedule:run
