# .root-fs/etc/s6-overlay/s6-rc.d/laravel/run
#!/usr/bin/execlineb -P

with-contenv
s6-setuidgid 10000:4

php82 /app/artisan octane:start --host=127.0.0.1 --port=8090
