# .root-fs/etc/s6-overlay/s6-rc.d/queue/run
#!/command/execlineb -P

with-contenv
s6-setuidgid www

php /app/artisan queue:work --sleep=3
