# .root-fs/etc/s6-overlay/s6-rc.d/reverb/run
#!/usr/bin/execlineb -P

with-contenv
s6-setuidgid www

php82 /app/artisan reverb:start --host=127.0.0.1 --port=8080
