# .root-fs/etc/s6-overlay/s6-rc.d/supercronic/run
#!/command/execlineb -P

/bin/supercronic /crontab
