# .root-fs/etc/s6-overlay/s6-rc.d/supercronic/run
#!/usr/bin/execlineb -P

/usr/bin/supercronic /crontab
