# .root-fs/etc/s6-overlay/s6-rc.d/supercronic/run
#!/usr/bin/execlineb -P

with-contenv
s6-setuidgid www

/usr/bin/supercronic /crontab
